# PiBuoy/sleepypid/sleepypid.service
[Unit]
Description=SleepyPi management

[Service]
Type=simple
Restart=always
ExecStart=/usr/local/bin/sleepypid.py --deepsleepvoltage 12.4 --shutdownvoltage 12.5 --fullvoltage 12.8 --log /var/log/sleepypid.log --sleepscript /scripts/shutdown.sh --startscript /scripts/startup.sh

[Install]
WantedBy=multi-user.target
