# PiBuoyV2/config/crontab
# every five minutes
*/5 * * * * /opt/AISonobuoy/PiBuoy/scripts/system_health.sh
# once an hour, one minute past the hour
1 * * * * /usr/bin/env bash /opt/AISonobuoy/PiBuoy/scripts/s3_prep.sh
